# roles/radioco/templates/supervisord.conf.in
[program:radioco_{{ branch }}]
autorestart=true
command={{ workdir }}/bin/gunicorn-with-settings -b {{ bind_address }} wrapper.wsgi
directory={{ workdir }}
redirect_stderr=true
user=radioco
stopasgroup=true
